# packages/app-build/resources/systemd/dassie-ipc.socket
[Unit]
Description=Dassie IPC Socket
Documentation=https://dassie.land/

[Socket]
ListenStream=/run/dassie.sock
SocketMode=0660
SocketUser=root
SocketGroup=dassie-users
Service=dassie.service

[Install]
WantedBy=sockets.target
