# packages/app-build/resources/systemd/dassie-https.socket
[Unit]
Description=Dassie Public HTTPS Socket
Documentation=https://dassie.land/

[Socket]
ListenStream=0.0.0.0:443
NoDelay=true
Service=dassie.service

[Install]
WantedBy=sockets.target
