# packages/app-build/resources/systemd/dassie-update.service
[Unit]
Description=Update Dassie to the latest version
Documentation=https://dassie.land/

[Service]
Type=oneshot
User=root
ExecStart=/opt/dassie/current/bin/dassie update
