# packages/app-build/resources/systemd/dassie-http.socket
[Unit]
Description=Dassie Public HTTP Socket
Documentation=https://dassie.land/

[Socket]
ListenStream=0.0.0.0:80
NoDelay=true
Service=dassie.service

[Install]
WantedBy=sockets.target
